# build/knxd/rootfs/etc/services.d/knxd/run
#!/bin/sh

/usr/local/bin/knxd /etc/knxd/knxd.ini
