# build/homebridge/rootfs/etc/services.d/homebridge/run
#!/bin/sh

homebridge -I 
#2>&1 >/var/log/homebridge.log
